# MMOCoreORB/docker/scripts/supervisord-core3.conf
[supervisord]
nodaemon=true

[program:core3]
user=root
command=/git/Core3/MMOCoreORB/bin/core3
directory=/git/Core3/MMOCoreORB/bin
